# ansible/config/supervisor/frontend.conf
[program:nuxt]
command=boot-nuxt
autostart=true
autorestart=true
stderr_logfile=/var/log/nuxt.err.log
stdout_logfile=/var/log/nuxt.out.log
stopasgroup=true
stopsignal=QUIT